# setup/setup/raspikey_usb.service
[Unit]
Description=raspikey_usb service

[Service]
ExecStart=/usr/sbin/raspikey_usb
User=root

[Install]
WantedBy=multi-user.target
